# supervisord.conf
[program:celeryworker]
command=celery worker -l info -A seoq.asynchronous
autorestart=true